# techjargon.celerybeat.conf
; ================================
;  celery beat supervisor example
; ================================

[program:celerybeat]
; Set full path to celery program if using virtualenv
command=celery -A techjargon beat -l info -S django

; remove the -A myapp argument if you aren't using an app instance

directory=/var/www/techjargon/techjargon-app
user=nobody
numprocs=1
stdout_logfile=/var/www/techjargon/techjargon-app/log/celery_beat.log
stderr_logfile=/var/www/techjargon/techjargon-app/log/celery_beat.log
logfile_maxbytes=50MB       ; maximum size of logfile before rotation
logfile_backups=10          ; number of backed up logfiles
loglevel=warn               ; info, debug, warn, trace
childlogdir=/var/www/techjargon/techjargon-app/log/celery_beat/            ; where child log files will live
autostart=true
autorestart=true
startsecs=10

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999
